# systemd/photon.service
[Unit]
Description=Photon Framework Application
Documentation=https://github.com/photon-framework/photon
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/photon
ExecStart=/opt/photon/app
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=photon

# Environment variables
Environment=APP_ENV=production
Environment=APP_PORT=8080
Environment=APP_HOST=0.0.0.0

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
ReadWritePaths=/opt/photon/data /opt/photon/logs

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
